# akademi-sihir.service
[Unit]
Description=Akademi Sihir Qithmir — AI RPG Game Server
After=network.target

[Service]
Type=simple
User=zevara
WorkingDirectory=/home/zevara/akademi-sihir
ExecStart=/home/zevara/akademi-sihir/run.sh 8000
Restart=always
RestartSec=10
Environment=OPENROUTER_MODEL=anthropic/claude-sonnet-4

[Install]
WantedBy=multi-user.target
